# deploy/systemd/openclaw-autopull.service
[Unit]
Description=OpenClaw auto pull from GitHub

[Service]
Type=oneshot
ExecStart=/usr/local/bin/openclaw-autopull.sh
